# docker-config/turkle.crontab
*/15 * * * * cd /opt/turkle && /usr/local/bin/python manage.py expire_assignments >> /var/log/cron.log 2>&1
