# docker-config/turkle.crontab
*/15 * * * * cd /opt/turkle && python manage.py expire_assignments >> /var/log/cron.log 2>&1
